# bash/profile.d/00-path.sh
export PATH="$HOME/.cask/bin:$HOME/bin:$HOME/.cabal/bin:/usr/local/bin:$PATH"
