# bash/profile.d/01-aliases.sh
#!/bin/bash

alias ls="ls --color=auto"
alias t=hodor
alias sagu="sudo apt-get update"
alias sagi="sudo apt-get install"
