# bash/profile.d/02-glob-star.sh
# If set, the pattern "**" used in a pathname expansion context will
# match all files and zero or more directories and subdirectories.

# Requires bash 4, at least.
shopt -s globstar
